# vbtspa/debian/postrm
#!/bin/sh
update-rc.d vbts_pa remove >/dev/null || true
# || exit $?

